# apps/minecraft/systemd/minecraft.target
[Unit]
Description=BrainCraft Minecraft Fleet

[Install]
WantedBy=default.target
