# apps/minecraft/quadlet/minecraft@.container
[Unit]
Description=Minecraft Java Edition - %i (BrainCraft)
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=minecraft-server.image
ContainerName=minecraft-%i
AutoUpdate=registry

Network=minecraft.network

Volume=%h/minecraft/%i:/data:Z
Volume=%h/quadlets/apps/minecraft/config:/config:ro,z

UserNS=keep-id:uid=1000,gid=1000
User=1000
Group=1000

HealthStartupCmd=pgrep -f java
HealthStartupInterval=10s
HealthStartupRetries=0
HealthStartupSuccess=1
HealthStartupTimeout=10s
HealthCmd=/usr/local/bin/mc-health
HealthInterval=30s
HealthRetries=5
HealthStartPeriod=0s
HealthTimeout=10s
HealthOnFailure=kill
Notify=healthy

StopTimeout=130

Memory=5g
PidsLimit=4096
Ulimit=nofile=65536:65536

NoNewPrivileges=true
DropCapability=all
Mask=/proc/acpi:/proc/kcore:/proc/keys:/proc/latency_stats:/proc/timer_list:/proc/timer_stats:/proc/sched_debug:/proc/scsi:/sys/firmware

HostName=minecraft-%i

LogDriver=journald
LogOpt=tag=minecraft-%i

Label=app=minecraft
Label=instance=%i
Label=managed-by=quadlet
Annotation=io.containers.description=BrainCraft_Minecraft_%i_Server

EnvironmentFile=%h/.config/containers/systemd/minecraft-secrets.env
EnvironmentFile=%h/.config/containers/systemd/minecraft-%i.env

Pull=missing
Retry=3
RetryDelay=5s

Timezone=America/Los_Angeles

# cgroup v2 direct writes
# Podman >= 5.9: use CgroupConf= key directly:
#   CgroupConf=memory.swap.max=0
#   CgroupConf=memory.high=4500000000
#   CgroupConf=memory.oom.group=1
# Podman 5.8: pass via --cgroup-conf flag
PodmanArgs=--cgroup-conf=memory.swap.max=0 --cgroup-conf=memory.high=4500000000 --cgroup-conf=memory.oom.group=1

[Quadlet]
DefaultDependencies=true

[Service]
Restart=on-failure
RestartSec=30
RestartSteps=4
RestartMaxDelaySec=300
OOMPolicy=stop
MemoryHigh=4500M
MemorySwapMax=0
TimeoutStartSec=900
TimeoutStopSec=160
SyslogIdentifier=minecraft-%i
CPUWeight=100

[Install]
WantedBy=default.target
