# apps/minecraft/quadlet/minecraft-proxy.container
[Unit]
Description=Minecraft Velocity Proxy (BrainCraft)
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=minecraft-proxy.image
ContainerName=minecraft-proxy
AutoUpdate=registry

Network=minecraft.network
PublishPort=25565:25565

Volume=%h/minecraft/proxy:/server:Z
Volume=%h/quadlets/apps/minecraft/config:/config:ro,z

UserNS=keep-id:uid=1000,gid=1000
User=1000
Group=1000

HealthStartupCmd=pgrep -f velocity
HealthStartupInterval=5s
HealthStartupRetries=0
HealthStartupSuccess=1
HealthStartupTimeout=5s
HealthCmd=mc-monitor status --host localhost --port 25565
HealthInterval=15s
HealthRetries=3
HealthStartPeriod=0s
HealthTimeout=5s
HealthOnFailure=kill
Notify=healthy

StopTimeout=30

Memory=2g
PidsLimit=512
Ulimit=nofile=65536:65536

NoNewPrivileges=true
DropCapability=all
Mask=/proc/acpi:/proc/kcore:/proc/keys:/proc/latency_stats:/proc/timer_list:/proc/timer_stats:/proc/sched_debug:/proc/scsi:/sys/firmware

HostName=minecraft-proxy

LogDriver=journald
LogOpt=tag=minecraft-proxy

Label=app=minecraft
Label=role=proxy
Label=managed-by=quadlet

EnvironmentFile=%h/.config/containers/systemd/minecraft-secrets.env
EnvironmentFile=%h/.config/containers/systemd/minecraft-proxy.env

Pull=missing
Retry=3
RetryDelay=5s

Timezone=America/Los_Angeles

# cgroup v2 direct writes
# Podman >= 5.9: use CgroupConf= key directly:
#   CgroupConf=memory.swap.max=0
#   CgroupConf=memory.high=900000000
# Podman 5.8: pass via --cgroup-conf flag
PodmanArgs=--cgroup-conf=memory.swap.max=0 --cgroup-conf=memory.high=1800000000 --cgroup-conf=memory.oom.group=1

[Quadlet]
DefaultDependencies=true

[Service]
Restart=on-failure
RestartSec=10
RestartSteps=4
RestartMaxDelaySec=120
OOMPolicy=kill
ManagedOOMMemoryPressure=kill
TimeoutStartSec=300
TimeoutStopSec=45
SyslogIdentifier=minecraft-proxy
CPUWeight=150
MemoryHigh=1800M
MemorySwapMax=0

# [Install] WantedBy= managed by deploy.sh via 20-fleet.conf drop-in.
# Proxy always gets WantedBy=minecraft.target.
